# contrib/reniced.service
[Unit]
Description=Reniced Service

[Service]
ExecStart=/usr/sbin/reniced
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
